# files/etc/uci-defaults/60-network
#!/bin/sh

uci set network.lan.ipv6='0'

uci commit network
